# third_party/process-exporter/packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable process-exporter.service
systemctl restart process-exporter.service
