# third_party/process-exporter/packaging/scripts/preremove.sh
systemctl stop process-exporter.service
systemctl disable process-exporter.service
