# build/node/hack/rpm/centos-7/files/usr/lib/systemd/system/acp-node.service
[Unit]
Description=goodrain acp-node
After=network.target
#Requires=etcd.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/goodrain/envs/acp-node.sh
PermissionsStartOnly=true
ExecStart=/usr/share/gr-acp-node/scripts/start-node.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target